# crontab
# START CRON JOB
*/15 * * * * python3 /app/app.py
# END CRON JOB